# standardtest/crontab
0 * * * * /src/publish_scheduled_pages.sh
